# hammerspoon-whisper/launchd/com.whisper.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whisper.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/whisper-server</string>
        <string>-m</string>
        <string>__HOME__/.ai/models/ggml-large-v3-turbo-q8_0.bin</string>
        <string>--port</string>
        <string>49440</string>
        <string>--vad</string>
        <string>--vad-model</string>
        <string>__HOME__/.ai/models/ggml-silero-v6.2.0.bin</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/var/log/whisper-server.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.local/var/log/whisper-server-error.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
